# idps-agent.service
[Unit]
Description=IDPS Agent - Intrusion Detection & Prevention System
After=network.target nginx.service fail2ban.service
Wants=network-online.target

[Service]
Type=simple
User=idps-agent
Group=idps-agent
WorkingDirectory=/opt/idps-agent
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Environment variables — fill in your actual values
Environment="TELEGRAM_BOT_TOKEN=your_token_here"
Environment="TELEGRAM_CHAT_ID=your_chat_id_here"
Environment="ANTHROPIC_API_KEY=your_anthropic_key_here"
Environment="AUTONOMOUS_MODE=false"
Environment="MONITORED_SERVICE=your-app-service-name"
Environment="API_BEARER_TOKEN=your_stats_endpoint_token"
Environment="ALLOWED_COUNTRIES=BR,US"
Environment="IDPS_PORT=3001"

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ReadWritePaths=/var/log/idps-agent /var/lib/idps-agent /etc/nginx/blocked-ips.conf
# Add only paths that exist on your server (missing paths cause startup failure).
# setup.sh will automatically uncomment paths for services present on this host.
ReadOnlyPaths=/var/log/auth.log
# Uncomment if nginx is installed
#ReadOnlyPaths=/var/log/nginx
# Uncomment if fail2ban is installed
#ReadOnlyPaths=/var/log/fail2ban.log
# Uncomment if ufw is installed
#ReadOnlyPaths=/var/log/ufw.log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=idps-agent

[Install]
WantedBy=multi-user.target
